# infrastructure/ansible/roles/api/templates/quiz-api.service.j2
[Unit]
Description={{ app_name }} API Server
After=network.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_group }}
WorkingDirectory={{ app_dir }}
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
Environment=NODE_ENV={{ node_env }}

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ app_dir }}

[Install]
WantedBy=multi-user.target